# goaccess/rootfs/etc/services.d/nginx/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the Nginx service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# shellcheck shell=bash

bashio::log.info "Stopping nginx service..."
